# meta-lxatac-bsp/recipes-core/lxatac-persistent-sysstat/files/var-log-sa.mount
[Mount]
Where=/var/log/sa
What=/srv/sysstat
Type=none
Options=bind
