# meta-lxatac-bsp/recipes-core/lxatac-persistent-sysstat/files/var-volatile-log-sa.mount
[Unit]
Before=local-fs.target

[Mount]
Where=/var/volatile/log/sa
What=/srv/sysstat
Type=none
Options=bind
