# meta-lxatac-bsp/recipes-core/lxatac-persistent-labgrid-cache/files/var-cache-labgrid.mount
[Unit]
Before=local-fs.target

[Mount]
Where=/var/cache/labgrid
What=/srv/cache/labgrid
Type=none
Options=bind
