# meta-lxatac-bsp/recipes-core/lxatac-persistent-labgrid-cache/files/var-cache-labgrid.mount
[Mount]
Where=/var/cache/labgrid
What=/srv/cache/labgrid
Type=none
Options=bind
