# meta-lxatac-bsp/recipes-core/lxatac-factory-data/files/lxatac-factory-data.service
[Unit]
Description=Gather and Distribute Factory Data
DefaultDependencies=no
After=systemd-remount-fs.service
Before=systemd-udevd.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/lxatac-factory-data

[Install]
WantedBy=sysinit.target
