# meta-lxatac-bsp/recipes-core/lxatac-persistent-journal/files/var-volatile-log-journal.mount
[Unit]
Before=local-fs.target

[Mount]
Where=/var/volatile/log/journal
What=/srv/journal
Type=none
Options=bind
