# meta-lxatac-bsp/recipes-core/lxatac-persistent-journal/files/var-log-journal.mount
[Mount]
Where=/var/log/journal
What=/srv/journal
Type=none
Options=bind
